# Armbian/24/armbian-bsp-cli-rk3128-box-current/rootfs/etc/update-motd.d/20-ip-info
#!/bin/bash
# Print Ethernet / Wi-Fi / WAN IPv4 and IPv6 lines for Armbian MOTD.
# Standalone: discovers addresses itself and prints formatted output.

set -o pipefail

safe_source() { [[ -f "$1" ]] && . "$1"; }

THIS_SCRIPT="ip-info"
MOTD_DISABLE=""

# Defaults (can be overridden by /etc/default/armbian-motd if present)
HIDE_IP_PATTERN='^dummy0|^lo|^docker|^hassio|^br-|^veth|^vnet|^virbr'
HIDE_LOCAL_IPV6="true"

# Optional overrides
safe_source /etc/default/armbian-motd
for f in ${MOTD_DISABLE}; do
	[[ "${f}" == "${THIS_SCRIPT}" ]] && exit 0
done

# one-time spacer after header (marks that something printed before Performance)
_motd_mark_after_header() {
  local base="/run/armbian-motd"
  mkdir -p "$base" 2>/dev/null || base="${XDG_RUNTIME_DIR:-/tmp}/armbian-motd"
  mkdir -p "$base" 2>/dev/null || return 0
  local stamp="${base}/.after_header.printed"
  [[ -e "$stamp" ]] || { echo ""; : >"$stamp"; }
}

get_wan_address()  { curl -4 -fsS --max-time 2 https://ipv4.whatismyip.akamai.com/ 2>/dev/null || true; }
get_wan6_address() { curl -6 -fsS --max-time 2 https://ipv6.whatismyip.akamai.com/ 2>/dev/null || true; }

classify_interface() {
  local intf="$1"

  if [[ -d "/sys/class/net/${intf}/wireless" || "${intf}" =~ ^wl ]]; then
    printf '%s\n' "wifi"
  else
    printf '%s\n' "ethernet"
  fi
}

get_lan_ips() {
  # echo "<eth4>|<wifi4>|<eth6>|<wifi6>"
  local eth4s=() wifi4s=() eth6s=() wifi6s=() intf class
  local add4=() add6=()

  for p in /sys/class/net/*; do
    intf="${p##*/}"
    [[ "$intf" =~ $HIDE_IP_PATTERN ]] && continue
    class="$(classify_interface "$intf")"

    # IPv4 (order preserved, unique)
    mapfile -t add4 < <(ip -4 -o addr show dev "$intf" 2>/dev/null \
      | awk '$3=="inet"{print $4}' | cut -d/ -f1 || true)
    for a in "${add4[@]}"; do
      case "$class" in
        wifi)
          [[ " ${wifi4s[*]} " == *" $a "* ]] || wifi4s+=("$a")
          ;;
        *)
          [[ " ${eth4s[*]} " == *" $a "* ]] || eth4s+=("$a")
          ;;
      esac
    done

    # IPv6 (prefer scope global when hiding link-local)
    if [[ "${HIDE_LOCAL_IPV6}" == "true" ]]; then
      mapfile -t add6 < <(ip -6 -o addr show dev "$intf" scope global 2>/dev/null \
        | awk '$3=="inet6"{print $4}' | cut -d/ -f1 || true)
    else
      mapfile -t add6 < <(ip -6 -o addr show dev "$intf" 2>/dev/null \
        | awk '$3=="inet6"{print $4}' | cut -d/ -f1 || true)
    fi
    for a in "${add6[@]}"; do
      case "$class" in
        wifi)
          [[ " ${wifi6s[*]} " == *" $a "* ]] || wifi6s+=("$a")
          ;;
        *)
          [[ " ${eth6s[*]} " == *" $a "* ]] || eth6s+=("$a")
          ;;
      esac
    done
  done

  printf '%s|%s|%s|%s\n' "${eth4s[*]}" "${wifi4s[*]}" "${eth6s[*]}" "${wifi6s[*]}"
}

print_ip_lines() {
  local eth4="$1" wifi4="$2" eth6="$3" wifi6="$4" wan4="$5" wan6="$6"
  local printed=0

  print_line() {
    local label="$1" value="$2" color="$3"
    [[ -z "$value" ]] && return 0
    (( printed )) || { _motd_mark_after_header; printed=1; }
    printf ' %-15s \x1B[%sm%s\x1B[0m\n' "${label}:" "${color}" "${value}"
  }

  # Drop WAN duplicates from local lists
  if [[ -n "$wan4" ]]; then
    eth4="$(printf '%s\n' ${eth4:-} | awk -v w="$wan4" 'NF && $0!=w' | paste -sd' ' -)"
    wifi4="$(printf '%s\n' ${wifi4:-} | awk -v w="$wan4" 'NF && $0!=w' | paste -sd' ' -)"
  fi
  if [[ -n "$wan6" ]]; then
    eth6="$(printf '%s\n' ${eth6:-} | awk -v w="$wan6" 'NF && $0!=w' | paste -sd' ' -)"
    wifi6="$(printf '%s\n' ${wifi6:-} | awk -v w="$wan6" 'NF && $0!=w' | paste -sd' ' -)"
  fi

  print_line "Eth IPv4" "${eth4// /, }" "92"
  print_line "WiFi IPv4" "${wifi4// /, }" "96"
  print_line "WAN IPv4" "${wan4}" "95"
  print_line "Eth IPv6" "${eth6// /, }" "92"
  print_line "WiFi IPv6" "${wifi6// /, }" "96"
  print_line "WAN IPv6" "${wan6}" "95"
}

main() {
  local combo eth4 wifi4 eth6 wifi6 wan4 wan6
  combo="$(get_lan_ips)"
  IFS='|' read -r eth4 wifi4 eth6 wifi6 <<< "$combo"
  wan4="$(get_wan_address)"
  wan6="$(get_wan6_address)"
  print_ip_lines "$eth4" "$wifi4" "$eth6" "$wifi6" "$wan4" "$wan6"
}

# If sourced, don't run main; if executed, run.
(return 0 2>/dev/null) || main
