# Armbian/24/armbian-bsp-cli-rk3128-box-current/rootfs/etc/cron.daily/armbian-quotes
#!/bin/bash

# Keep local MOTD tips fully package-managed for RK3128 images.
# This package must not pull quotes from Armbian infrastructure.
exit 0
